# contrib/idiotic.service
[Service]
Type=simple
ExecStart=/usr/bin/idiotic -b /etc/idiotic
